# deploy/launchd/com.artjeck.ifta-web.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!--
  IFTA web intake service. Run by launchd at login and kept alive.
  Install with deploy/install.sh (substitutes __PROJECT_ROOT__).

  Logs land in __PROJECT_ROOT__/logs/web.{out,err}.log.
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.artjeck.ifta-web</string>

    <key>WorkingDirectory</key>
    <string>__PROJECT_ROOT__</string>

    <key>ProgramArguments</key>
    <array>
        <string>__PROJECT_ROOT__/.venv/bin/ifta</string>
        <string>web</string>
        <string>--host</string>
        <string>127.0.0.1</string>
        <string>--port</string>
        <string>8000</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <!-- Don't thrash on crash loops; back off if the process dies inside 10s. -->
    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>StandardOutPath</key>
    <string>__PROJECT_ROOT__/logs/web.out.log</string>

    <key>StandardErrorPath</key>
    <string>__PROJECT_ROOT__/logs/web.err.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
